# parcial1/cookbooks/flask/files/default/flask_service.conf
description     "my job"
start on startup
task
exec /home/flask/smartlabs/script.sh